# shinken/supervisor/conf.d/snmpd.conf
[program:snmpd]
command=/usr/sbin/snmpd -f /etc/snmp/snmpd.conf -p /var/run/snmpd.pid
process_name=%(program_name)s
autostart=true
autorestart=true
user=root
directory=/
priority=958
stdout_logfile=/var/log/supervisor_%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor_%(program_name)s-stderr.log

